# systemd/quectool.service
[Unit]
Description=Quectool
After=ql-netd.service

[Service]
ExecStart=/usrdata/quectool server
Restart=always
RestartSec=1s
Environment="SERVER_AUTH_REALM=Quectool"

[Install]
WantedBy=multi-user.target
